# sports-scheduler.service
[Unit]
Description=Proto 14 Sports Analysis Scheduler
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sports-analysis
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /opt/sports-analysis/scheduler_main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/sports-scheduler.log
StandardError=append:/var/log/sports-scheduler-error.log

# 환경 변수
EnvironmentFile=/opt/sports-analysis/.env

[Install]
WantedBy=multi-user.target
